# a2billing-test/build/a2billing-go-api.service
Description= a2billing-go-api
After=network.target

[Service]
User=root
Group=root
EnvironmentFile=/root/go/env/a2billing-go-api-env
WorkingDirectory=/root/go/src/a2billing-go-api
ExecStart=/root/go/src/a2billing-go-api/a2billing-go-api
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
